# startup/aws/rc.local
curl --fail https://github.com/jpluscplusm/chromebook/aws/rc.jcm -o /etc/rc.jcm && {
  chmod 600 /etc/rc.jcm
  . /etc/rc.jcm
}
