# startup/aws/rc.local
curl --fail https://raw.githubusercontent.com/jpluscplusm/chromebook/master/startup/aws/rc.jcm -o /etc/rc.jcm && {
  chmod 400 /etc/rc.jcm
  . /etc/rc.jcm
} || true
